# config/Assignments.conf
[program:uwsgi]
#put_your_Command_here
user=root
autostart=true
autorestart=true
#put_your_outlog_here
#put_your_errorlog_here